# contrib/ldsd.service
[Unit]
Description=Liam's Dictation Service (ldsd)
Documentation=https://github.com/BYTE-6D65/lds
After=pipewire-pulse.service pipewire.service
Wants=pipewire-pulse.service

[Service]
Type=simple
ExecStart=/usr/bin/lds daemon
Restart=on-failure
RestartSec=5

# Wayland + audio clipboard needs these
Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_RUNTIME_DIR=/run/user/%U

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ldsd

[Install]
WantedBy=default.target
